# ext/templates/puppetdb.service.erb
[Unit]
Description=<%= @name %> Daemon
After=syslog.target network.target

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/<%= @name %>
User=${USER}
PIDFile=/var/run/${USER}/${USER}.pid
ExecStart=${JAVA_BIN} \
          ${JAVA_ARGS} \
          -jar ${INSTALL_DIR}/puppetdb.jar \
          services -c ${CONFG} $@

ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target
