# ext/templates/puppetdb.service.erb
[Unit]
Description=<%= @name %> Daemon
After=syslog.target network.target

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/<%= @name %>
User=<%= @name %>
PIDFile=/var/run/<%= @name %>/<%= @name %>.pid
<%- if @pe -%>
ExecStart=<%= @java_bin %> \
          $JAVA_ARGS \
          -jar ${INSTALL_DIR}/puppetdb.jar \
          services -c ${CONFIG} $@
<%- else -%>
ExecStart=/usr/lib/jvm/jre-1.7.0-openjdk/bin/java \
          $JAVA_ARGS \
          -jar ${INSTALL_DIR}/puppetdb.jar \
          services -c ${CONFIG} $@
<%- end -%>

ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target
