# ext/templates/puppetdb.service.erb
[Unit]
Description=<%= @name %> Daemon
After=syslog.target network.target

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/<%= @name %>
User=<%= @name %>
ExecStart=/usr/lib/jvm/jre-1.7.0-openjdk/bin/java \
          $JAVA_ARGS \
          -cp ${INSTALL_DIR}/puppetdb.jar \
          clojure.main -m puppetlabs.puppetdb.core \
          services -c $CONFIG $@

[Install]
WantedBy=multi-user.target
