# ext/templates/deb/postinst.erb
#!/bin/bash

<% if @pe -%>
export PATH=/opt/puppet/bin:$PATH
<% end -%>

# Setup the SSL stuff

if [ ! -f "<%= @config_dir -%>/../ssl/puppetdb_keystore_pw.txt" ] ; then
  <%= @sbin_dir -%>/puppetdb-ssl-setup
fi

<%= ERB.new(File.read("ext/templates/directory_perms.erb")).result %>


set -e
if [ -x "/etc/init.d/puppetdb" ]; then
    if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
        invoke-rc.d puppetdb start || exit $?
    else
        /etc/init.d/puppetdb start || exit $?
    fi
fi
